# openwrt/zapret2-orchestra/files/etc/init.d/zapret2-orchestra-learner
#!/bin/sh /etc/rc.common

# zapret2-orchestra-learner — procd-managed closed-loop learning daemon.
#
# Starts the ucode learner that tails /tmp/zapret2-orchestra/events.ndjson,
# persists learned/blocked/learner-state, and reloads nfqws2 (debounced) when a
# lock or blocked-strategy changes.  START after the zapret2 service so the
# events sink and datapath exist before the learner begins; STOP cleanly
# terminates the daemon (the durable cursor makes restart safe).
#
# The reload the learner requests is /etc/init.d/zapret2 restart, owned HERE by
# the procd layer (the learner daemon itself never references /etc/init.d from
# the packet path — it sets ORCHESTRA_RELOAD_CMD to this service's reload
# helper, keeping the packet-path-writes-no-JSON / no-init.d-from-apply
# separation intact).

START=22
USE_PROCD=1

LEARNER_BIN=/usr/sbin/zapret2-orchestra-learner
# The reload helper this service owns: restart the zapret2 datapath.  The
# learner daemon uses ORCHESTRA_RELOAD_CMD to invoke it (debounced), so the
# /etc/init.d reference lives in the procd layer, not in apply.uc or learner.uc.
RELOAD_HELPER=/usr/sbin/zapret2-orchestra-reload

ORCHESTRA_RUNTIME_DIR=/tmp/zapret2-orchestra

# Ensure the volatile Orchestra runtime dir exists with secure ownership so the
# unprivileged nfqws2 user (daemon) can append events.ndjson, while root-owned
# preload/manifest/whitelist files stay read-only to it.  /tmp is cleared on
# reboot, so this recreates the dir when the learner (re)starts.  Rejects a
# symlinked path, never makes the dir world-writable, never touches the
# persistent /etc/zapret2-orchestra state dir, and never truncates
# events.ndjson (unread events survive a restart).  Best-effort: a chown
# failure is logged but does not block the learner (the boot hook is the
# primary guarantor; this is belt-and-suspenders for standalone restarts).
ensure_runtime_dir() {
	[ -L "$ORCHESTRA_RUNTIME_DIR" ] && {
		logger -t zapret2-orchestra "runtime dir $ORCHESTRA_RUNTIME_DIR is a symlink; refusing" 2>/dev/null
		return 1
	}
	mkdir -p "$ORCHESTRA_RUNTIME_DIR" || return 1
	chown daemon:daemon "$ORCHESTRA_RUNTIME_DIR" 2>/dev/null || \
		logger -t zapret2-orchestra "warning: cannot chown $ORCHESTRA_RUNTIME_DIR daemon:daemon" 2>/dev/null
	chmod 0770 "$ORCHESTRA_RUNTIME_DIR"
}

start_service() {
	# The events sink dir must exist (daemon:daemon 0770) before the learner
	# tails it and before nfqws2 appends events.ndjson.  Best-effort: a failure
	# is logged but does not block the learner daemon.
	ensure_runtime_dir || true
	# The lock dir is under /var/lock (volatile); create it so the learner's
	# mkdir-lock can write a pid file.
	mkdir -p /var/lock
	procd_open_instance
	procd_set_param command "$LEARNER_BIN" run
	# The learner reloads nfqws2 via this helper (procd-owned).
	procd_set_param env ORCHESTRA_RELOAD_CMD="$RELOAD_HELPER"
	procd_set_param respawn 3600 5 30
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

reload_service() {
	# A manual `/etc/init.d/zapret2-orchestra-learner reload` runs the learner
	# once (catch up on any events emitted while the daemon was down) and then
	# reloads the datapath.
	"$LEARNER_BIN" process-once
	"$RELOAD_HELPER"
}

stop_service() {
	# procd handles SIGTERM; the learner's durable cursor makes a clean stop
	# safe (no state loss).
	:
}
